# services/processing-worker/crontab
# Run cleanup job daily at 2 AM UTC
# Format: minute hour day month weekday command
0 2 * * * /usr/local/bin/cleanup >> /proc/1/fd/1 2>&1
